# Solutions/Tier0/roles/filesystem_setup/templates/hammerspace-remount.timer.j2
# Hammerspace Auto-Remount Watchdog Timer
# Periodically triggers the remount check service
# Generated by Ansible - do not edit manually

[Unit]
Description=Hammerspace Auto-Remount Watchdog Timer

[Timer]
OnBootSec=30
OnUnitActiveSec={{ hammerspace_remount_watchdog_interval | default('1min') }}
Unit=hammerspace-remount.service

[Install]
WantedBy=timers.target
