# Solutions/Tier0/roles/filesystem_setup/templates/hammerspace-guard.service.j2
# Hammerspace Mount Guard Service for {{ mount_path }}
# Prevents accidental unmount by keeping a process with cwd on the mount point
# Generated by Ansible - do not edit manually

[Unit]
Description=Hammerspace {{ mount_name }} Mount Guard (busy-lock)
After={{ mount_unit }}.mount
Requires={{ mount_unit }}.mount
BindsTo={{ mount_unit }}.mount

[Service]
Type=simple
ExecStart=/bin/bash -lc 'cd {{ mount_path }} && exec sleep infinity'
Restart=always
RestartSec=1
KillMode=process
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
OOMScoreAdjust=-1000
RefuseManualStop=yes

[Install]
WantedBy=hammerspace-guards.target
