# roles/filesystem_setup/templates/hammerspace-guards.target.j2
# Hammerspace Mount Guards Target
# Manages all mount guard services together
# Generated by Ansible - do not edit manually

[Unit]
Description=Hammerspace Mount Guards
Documentation=man:systemd.special(7)
After=local-fs.target

[Install]
WantedBy=multi-user.target
