# Tier0/roles/filesystem_setup/templates/hammerspace-remount.service.j2
# Hammerspace Auto-Remount Watchdog Service
# Checks and remounts Hammerspace volumes if accidentally unmounted
# Generated by Ansible - do not edit manually

[Unit]
Description=Hammerspace Auto-Remount Watchdog
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hammerspace-remount-check.sh

[Install]
WantedBy=multi-user.target
